# templates/cni-up.timer.epp
[Unit]
Description=Keep VXLAN CNI up

[Timer]
OnBootSec=2m
OnUnitActiveSec=<%= $keep_vxlan_up_cron_interval %>m

[Install]
WantedBy=timers.target
