# templates/cni-up.service.epp
[Unit]
Description=VXLAN CNI service for timer execution
Documentation=https://github.com/GEANT/nomad_cni
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment="STARTED_BY_SYSTEMD=yes"
User=root
Group=root
TimeoutStartSec=0
ExecStart=/usr/local/bin/cni-vxlan-wizard.sh --status up --name all<% if $ingress !~ Undef { %> --vip<% } %>
SyslogIdentifier=CNI-VXLAN-wizard

[Install]
WantedBy=multi-user.target
