# files/cni-purge.timer
[Unit]
Description=Keep VXLAN CNI up

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target