# files/cni-purge.timer
[Unit]
Description=Purge VXLAN interfaces without a matching configuration

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target