# files/rc.local
#!/bin/sh
#
# This script is managed by the puppet module nomad_cni - do not modify
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

/usr/local/bin/vxlan-configurator.sh --force --all

exit 0